# misc/testing/DB_scripts/populate_anidb.cron
#!/bin/bash
export PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/lib/jvm/java-7-oracle/bin
php5 /websites/nZEDb/misc/testing/DB_scripts/populate_anidb_with_proxy.php 30 >> /websites/nZEDb/logs/cron.log
