# misc/testing/Dev_testing/populate_anidb.cron
#!/usr/bin/env bash
# This file should be copied/dropped into your cron path

export DB_PATH="/var/www/nZEDb/misc/testing/DB_scripts"
export LOG_PATH="/var/www/nZEDb/logs"
export PATH=$PATH

command -v php5 >/dev/null 2>&1 && export PHP=`command -v php5` || { export PHP=`command -v php`; }

$PHP ${DB_PATH}/populate_anidb.php 30 cron >> ${LOG_PATH}/cron.log
